# init/dnsperf_daemon
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="DNS Performance Daemon"
description="Monitors DNS performance by testing response times periodically"

command="/usr/local/bin/dns_perf_backend.sh"
command_user="root"
pidfile="/var/run/dnsperf_daemon.pid"
command_background="yes"

depend() {
    need net
    after logger
}

start_pre() {
    # Ensure working directory exists
    checkpath --directory --owner root:root --mode 0755 /var/lib/dnsperf_daemon

    # Check if daemon script exists and is executable
    if [ ! -x "${command}" ]; then
        eerror "Daemon script not found or not executable: ${command}"
        return 1
    fi

    # Check dependencies
    if ! command -v dnsperf >/dev/null 2>&1; then
        eerror "dnsperf is required but not installed"
        return 1
    fi

    if ! command -v wget >/dev/null 2>&1; then
        eerror "wget is required but not installed"
        return 1
    fi

    if ! command -v unzip >/dev/null 2>&1; then
        eerror "unzip is required but not installed"
        return 1
    fi

    einfo "Starting ${name}..."
    return 0
}

start_post() {
    einfo "${name} started successfully"
}

stop_post() {
    einfo "${name} stopped"
}

reload() {
    einfo "Reloading ${name} configuration..."
    if [ -f "${pidfile}" ]; then
        pid=$(cat "${pidfile}")
        if kill -HUP "${pid}" 2>/dev/null; then
            einfo "${name} configuration reloaded successfully"
            return 0
        else
            eerror "Failed to send reload signal to ${name}"
            return 1
        fi
    else
        eerror "${name} is not running"
        return 1
    fi
}
