# init/dnsperf_daemon.service
[Unit]
Description=DNS Performance Monitoring Daemon
Documentation=man:dnsperf(1)
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/dns_perf_backend.sh daemon
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/dnsperf_daemon /var/log /var/run
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Working directory and environment
WorkingDirectory=/var/lib/dnsperf_daemon
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Process management
PIDFile=/var/run/dnsperf_daemon.pid
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
